# src/main/resources/hdfs/hbase-pure/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:hdfs]
command=/bin/bash -c "/hadoop/hadoop-2.10.2/sbin/start-all.sh"

[program:hdfsinit]
command=/bin/bash -c "/hdfs-init.sh"
stopasgroup=true
stopsignal=KILL
startretries=10

[group:upfuzz_hdfs]
programs=hdfs,hdfsinit
